# packaging/macos/scripts/postinstall
#!/bin/bash

set -e

echo "SafeChain Ultimate Installer - Post-installation"

INSTALL_DIR="/Library/Application Support/AikidoSecurity/SafeChainUltimate"
LOGS_DIR="/Library/Logs/AikidoSecurity/SafeChainUltimate"
LAUNCHDAEMON_PLIST="/Library/LaunchDaemons/com.aikidosecurity.safechainultimate.plist"

# Ensure logs directory exists and has proper permissions
echo "Setting up logs directory..."
mkdir -p "$LOGS_DIR"
chmod 755 "$LOGS_DIR"

# Set proper permissions on binaries
echo "Setting permissions on binaries..."
chmod 755 "$INSTALL_DIR/bin/safechain-ultimate"
chmod 755 "$INSTALL_DIR/bin/safechain-ultimate-ui"
chmod 755 "$INSTALL_DIR/bin/safechain-proxy"

# Set proper ownership (root:wheel for system daemon)
echo "Setting ownership..."
chown -R root:wheel "$INSTALL_DIR"
chown -R root:wheel "$LOGS_DIR"

if [ -n "$AIKIDO_TOKEN" ]; then
    echo "Storing Aikido token..."
    TOKEN_FILE="$INSTALL_DIR/.token"
    printf '%s' "$AIKIDO_TOKEN" > "$TOKEN_FILE"
    chmod 600 "$TOKEN_FILE"
    chown root:wheel "$TOKEN_FILE"
fi

# Load and start the LaunchDaemon
if [ -f "$LAUNCHDAEMON_PLIST" ]; then
    echo "Loading LaunchDaemon..."
    launchctl bootstrap system "$LAUNCHDAEMON_PLIST" 2>/dev/null || {
        # Fallback to launchctl load (old way of starting the daemon on older macOS versions)
        echo "LaunchDaemon bootstrap failed, trying load..."
        launchctl load "$LAUNCHDAEMON_PLIST" || true
    }
    echo "SafeChain Ultimate daemon installed and started"
else
    echo "Warning: LaunchDaemon plist not found at $LAUNCHDAEMON_PLIST"
fi

echo "Post-installation completed"
echo ""
echo "SafeChain Ultimate has been installed successfully!"
echo "  Binaries: $INSTALL_DIR/bin"
echo "  Logs: $LOGS_DIR"
echo ""
echo "The agent is now running as a system daemon."

exit 0
